# dotfiles/swaybg-hourly.service
[Unit]
Description=Run the swaybg-hourly service

[Service]
ExecStart=/usr/local/bin/swaybg-hourly.sh
